# infra/highhopes-menu-stage.service
[Unit]
Description=High Hopes Menu API (staging)
After=network.target

[Service]
User=root
WorkingDirectory=/home/highhopes/highhopes-menu-stage/backend
EnvironmentFile=/home/highhopes/highhopes-menu-stage/backend/.env
ExecStart=/home/highhopes/highhopes-menu-stage/backend/venv/bin/waitress-serve --host 127.0.0.1 --port 5002 app:app
Restart=always
RestartSec=5
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
