# infra/highhopes-menu-stage.service
[Unit]
Description=High Hopes Menu API (staging)
After=network.target

[Service]
User=root
WorkingDirectory=/home/highhopes/highhopes-menu-stage/backend
EnvironmentFile=/home/highhopes/highhopes-menu-stage/backend/.env
ExecStartPre=/bin/rm -f /home/highhopes/highhopes-menu-stage/backend/gunicorn.ctl
ExecStart=/home/highhopes/highhopes-menu-stage/backend/venv/bin/gunicorn -w 1 --timeout 30 -b 127.0.0.1:5002 app:app
Restart=always
RestartSec=5
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
